# init/isletd.conf
description "Isletd"
author "Jon Schipp"
start on runlevel [2345]

exec /opt/islet/bin/isletd
