# cron/islet.crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command

# Don't run at the same time, causes problems
@daily   	 root    LOCATION/remove_old_containers
* * * * *        root    LOCATION/remove_old_users
* * * * *        root    LOCATION/port_forward
@hourly   	 root    LOCATION/disk_limit
@hourly          root    LOCATION/upload_container
